# server/hardening/os-isolation/multi-tenant@.service
[Unit]
Description=Multi-Tenant Node.js Service (Tenant: %i)
After=network.target mongodb.service redis.service
# Fail hard if tenant configuration is missing
ConditionPathExists=/etc/multi-tenant/tenants/%i.env

[Service]
Type=simple
# Security: Execute as dedicated tenant user/group
User=tenant-%i
Group=tenant-%i
UMask=0077

# Filesystem Isolation
WorkingDirectory=/opt/multi-tenant/current
EnvironmentFile=/etc/multi-tenant/tenants/%i.env
ExecStart=/usr/bin/node src/app.js
Restart=always
RestartSec=10

# Hardening Limits
# Prevent the process from gaining new privileges (e.g. via setuid binaries)
NoNewPrivileges=yes

# Isolation: Protect system directories
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Restricted access to /proc and /sys
ProcSubset=pid
ProtectProc=invisible

# Filesystem access control
# Only allow writing to the dedicated tenant data and log directories
ReadOnlyPaths=/
ReadWritePaths=/var/lib/multi-tenant/%i /var/log/multi-tenant/%i /tmp /var/tmp

# Disable all capabilities (Node.js doesn't need them for standard web apps)
CapabilityBoundingSet=
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
